# deploy/systemd/cloudflared-tunnel.service
[Unit]
Description=Cloudflared Quick Tunnel for HomeVault
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/server
Environment=PYTHONPATH=%h/server
ExecStart=%h/server/venv/bin/python %h/server/scripts/cloudflared_tunnel.py
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
